# systemd/daily.service
[Unit]
Description=Daily exec of shell script

[Service]
Type=simple
ExecStart=/path/to/script.sh

[Install]
WantedBy=multi-user.target
